# extras/backup-lancer.timer
[Unit]
Description=Back up Mr Lancer's database daily

[Timer]
OnCalendar=daily

[Install]
WantedBy=timers.target
